# dockerizedSauna/WellnessUSA-Solfeggio/sound-player-backend.service
[Unit]
Description=JJ9 Hz Sound Player Backend (FastAPI)
After=network.target sound.target alsa-state.service
Wants=network-online.target
Requires=alsa-state.service

[Service]
Type=simple
User=admin
SupplementaryGroups=audio
WorkingDirectory=/home/admin/Python
Environment="PATH=/home/admin/.local/share/virtualenvs/python-*/bin:/usr/local/bin:/usr/bin:/bin"
Environment="ALSA_CARD=0"
Environment="SDL_AUDIODRIVER=alsa"
ExecStart=/home/admin/Python/start-backend.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
DeviceAllow=/dev/snd/* rw
PrivateDevices=no

[Install]
WantedBy=multi-user.target

